# systemd/tor-ext-tun2socks.service.in
[Unit]
Description=tor-ext transparent proxy (tun2socks -> distro tor SOCKS5)
Documentation=https://github.com/vipinus/gnome-shell-extension-tor
# Bind to BOTH possible distro unit names — only one will exist on a given
# host. systemd ignores BindsTo entries that resolve to "not-found" with a
# warning rather than failing, so this stays dual-distro friendly without
# needing install-time templating.
BindsTo=tor.service
BindsTo=tor@default.service
After=tor.service tor@default.service network-online.target
Wants=network-online.target

[Service]
Type=simple
# v0.6.0: runs as root. The previous _tor-ext-uid sandbox split has been
# removed because the routing helper invoked from ExecStartPost needs nft
# NAT tables, ip rule, and per-iface sysctls — all root-only operations
# that can't be granted via CAP_NET_ADMIN alone. Splitting privileges
# across two units to save one dropped-cap process is not worth the
# cross-unit failure modes.

Environment=TUN_DEV=@TUN_DEV@
Environment=TUN_ADDR=@TUN_ADDR@
Environment=SOCKS_PORT=@SOCKS_PORT@
Environment=DNS_PORT=@DNS_PORT@

# Pre-create the TUN. Owned by root (matches unit privilege).
ExecStartPre=@IP_BIN@ tuntap add dev ${TUN_DEV} mode tun
ExecStartPre=@IP_BIN@ addr add ${TUN_ADDR} dev ${TUN_DEV}
ExecStartPre=@IP_BIN@ link set dev ${TUN_DEV} up mtu 1500

# xjasonlyu/tun2socks CLI. systemd resolves the executable path BEFORE
# env-var substitution, so the path must be baked in by the installer
# (not ${TUN2SOCKS}).
ExecStart=@TUN2SOCKS@ -device ${TUN_DEV} -proxy socks5://127.0.0.1:${SOCKS_PORT} -loglevel warning

# Policy routing + DNS redirect (TOR_UID detected at runtime from the
# distro tor user — debian-tor / tor / _tor).
ExecStartPost=/usr/local/libexec/tor-ext/tor-ext-routing up
ExecStopPost=/usr/local/libexec/tor-ext/tor-ext-routing down
ExecStopPost=-@IP_BIN@ tuntap del dev ${TUN_DEV} mode tun

# /run/tor-ext for the routing helper to stash IPv6-default state across
# up/down. Auto-cleaned by systemd on stop.
RuntimeDirectory=tor-ext
RuntimeDirectoryMode=0755

Restart=on-failure
RestartSec=3s

[Install]
WantedBy=multi-user.target
